# root/etc/services.d/discord/run
#!/usr/bin/execlineb -S0

exec php /var/www/Dramiel.php